# root/etc/services.d/sslsplit/run
#!/usr/bin/with-contenv bash

sslsplit -D -j /data/sslsplit/ -k /key/ca.key -c /key/ca.crt https 0.0.0.0 443 http 0.0.0.0 80

sleep infinity &

wait
